# rpm/systemd/diamond.service
[Unit]
Description=diamond - A system statistics collector for graphite

[Service]
ExecStart=/usr/bin/python /usr/bin/diamond --foreground
Restart=on-abort

[Install]
WantedBy=multi-user.target
